# deploy/supervisord.conf
# Two coupled processes in one container: the unshackle engine (REST API on
# 127.0.0.1:8786) and the Telegram bot that drives it. Logs go to the container stdout.
[program:engine]
command=unshackle serve --no-key --api-only --port 8786
autorestart=true
startsecs=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bot]
command=python -m src.bot
directory=/app
autorestart=true
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
